# deploy/systemd/exchange-events-web.service
# Web service (API + dashboard). Copy to /etc/systemd/system/, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now exchange-events-web
#
# Runs as its own unit, own working directory, own user -- deliberately
# independent of any other service/dashboard on the same host. A restart or
# crash here never touches other units.
[Unit]
Description=Exchange Events Pipeline - web (API + dashboard)
After=network.target

[Service]
Type=notice
User=exchange-events
Group=exchange-events
WorkingDirectory=/opt/exchange-events
EnvironmentFile=/opt/exchange-events/.env
ExecStart=/opt/exchange-events/.venv/bin/gunicorn wsgi:app \
    --bind 0.0.0.0:8080 \
    --workers 4 \
    --access-logfile /var/log/exchange-events/access.log \
    --error-logfile /var/log/exchange-events/error.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
